# systemd/xhs-watcher.service
[Unit]
Description=xhs-watcher — XHS keyword scrape → LLM verdict → terminal+TG
Documentation=file:///home/damian/xhs-watcher/DESIGN.md

[Service]
Type=oneshot
WorkingDirectory=/home/damian/xhs-watcher
ExecStart=/home/damian/xhs-watcher/run-pipeline.sh
StandardOutput=append:/tmp/xhs-watcher.log
StandardError=append:/tmp/xhs-watcher.log

[Install]
WantedBy=default.target
